# idle-daemon/assets/idle-daemon.service
[Unit]
Description=IdleScreen screensaver daemon
Documentation=https://github.com/idlescreen/idle-core
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=dbus
BusName=io.github.ubermetroid.trance
ExecStartPre=/usr/bin/mkdir -p %h/.config/idle
ExecStart=/usr/bin/idle-daemon daemon
# D-Bus well-known name remains historical for client ABI.
Restart=on-failure
RestartSec=5
TimeoutStopSec=10
KillMode=mixed

NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=%t -%h/.config/idle -%h/.config/trance
PrivateTmp=true
ProtectHome=read-only
MemoryDenyWriteExecute=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
LockPersonality=true
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallArchitectures=native

[Install]
WantedBy=graphical-session.target
Alias=idlescreen-daemon.service
Alias=idle-daemon.service
